# image/etc/supervisord.d/crond.ini
[program:crond]
command=/sbin/crond -n
redirect_stderr=true
