# image/etc/supervisord.d/crond.ini
[program:crond]
command=/sbin/crond -n
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=3
