# image/etc/supervisord.d/elasticsearch.ini
[program:elasticsearch]
command=/usr/share/elasticsearch/bin/elasticsearch
    -Des.default.path.data=/data/elasticsearch
    -Des.cluster.name=docker-suricata-elk
    -Des.default.path.logs=/var/log/elasticsearch
    -Des.default.config=/etc/elasticsearch/elasticsearch.yml
user=user
redirect_stderr=true
