# image/etc/supervisord.d/logstash.ini
[program:logstash]
command=/opt/logstash/bin/logstash --config /etc/logstash/conf.d
redirect_stderr=true
autostart=yes
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=3
user=user
