# image/etc/supervisord.d/logstash.ini
[program:logstash]
command=/opt/logstash/bin/logstash --config /etc/logstash/conf.d
redirect_stderr=true
autostart=yes
