# image/etc/supervisord.d/nginx.ini
[program:nginx]
command=/sbin/nginx -g "daemon off;"
redirect_stderr=true
