# image/etc/supervisord.d/nginx.ini
[program:nginx]
command=/sbin/nginx -g "daemon off;"
redirect_stderr=true
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=3
